# setup/pi_alarm_daemon.orig
#!/bin/sh

### BEGIN INIT INFO
# Provides:          pi_alarm_daemon
# Required-Start:    autofs $all
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: pi_alarm daemon
# Description:       pi_alarm daemon checking alarms and actions
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/opt/pi_alarm/pi_alarm_daemon.py
DAEMON_NAME=pi_alarm_daemon
DESC="Raspberry Pi NTP alarm clock daemon"

# This next line determines what user the script runs as.
# Root generally not recommended but necessary if you are using the Raspberry Pi GPIO from Python.
DAEMON_USER=root

test -x $DAEMON || exit 0

. /lib/lsb/init-functions

set -e

case "$1" in
  start)
    echo -n "Starting $DESC: "
    start-stop-daemon --start --quiet --pidfile /var/run/$DAEMON_NAME.pid --background --make-pidfile --exec $DAEMON
    echo "$DAEMON_NAME."
    ;;
  stop)
    echo -n "Stopping $DESC: "
    start-stop-daemon --stop --oknodo --quiet --pidfile /var/run/$DAEMON_NAME.pid
    echo "$DAEMON_NAME."
    ;;
  restart)
    $0 stop
    sleep 1
    $0 start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}" >&2
    exit 1
    ;;

esac

exit 0
